# Docker/runit/scraper-api/run
#!/bin/sh

cd /api
exec 2>&1
exec /usr/bin/npm run start
