# Docker/runit/xorg/run
#!/bin/sh

exec 2>&1
exec /usr/bin/Xorg -noreset +extension GLX +extension RANDR +extension RENDER -config ./xorg.conf 
